# ops/autostart/macos/com.agiletec.airis-mcp-gateway.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agiletec.airis-mcp-gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{DOCKER_COMPOSE_PATH}}</string>
        <string>-f</string>
        <string>{{REPO_ROOT}}/docker-compose.yml</string>
        <string>up</string>
        <string>-d</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/airis-mcp-gateway.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/airis-mcp-gateway.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:{{HOME}}/.orbstack/bin</string>
    </dict>
</dict>
</plist>
